# etc/crontab
#
# $Id: crontab,v 1.2 2007/10/05 12:43:44 az Exp az $
#

# m  h  dom mon dow command
0/10 *  *   *   *   ${HOME}/bin/ibm-mydyn.pl >${HOME}/tmp/ibm-mydyn.log 2>&1
